# systemd/filmbot-daily-report.timer
[Unit]
Description=Filmbot Daily Health Report Timer
Requires=filmbot-daily-report.service

[Timer]
# Run daily at 8:00 AM
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target

